# systemd/dnscrypt-proxy.socket
[Unit]
Description=dnscrypt-proxy listening socket

[Socket]
ListenStream=53
ListenDatagram=53
NoDelay=true
DeferAcceptSec=1

[Install]
WantedBy=sockets.target
